# nanoclaw.service
[Unit]
Description=NanoClaw - Personal Claude Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/behole/nanoclaw
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=append:/home/behole/nanoclaw/logs/nanoclaw.log
StandardError=append:/home/behole/nanoclaw/logs/nanoclaw.error.log

# Environment
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin

[Install]
WantedBy=default.target
